# systemd/LidarCounter.service
[Unit]
Description=Lidar Car Counter
After=network.target

[Service]
Type=simple

# --- Absolute paths (do not change) ---
WorkingDirectory=/home/admin/LidarCounter
ExecStart=/home/admin/LidarCounter/venv/bin/python /home/admin/LidarCounter/app.py

Restart=always
RestartSec=5

# Environment
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1

# User permissions (required for UART)
User=admin
Group=admin

[Install]
WantedBy=multi-user.target
